# software/chess-board/deploy/systemd/chessboard-dpad.service
[Unit]
Description=Physical Chessboard D-Pad Keyboard Bridge
After=chessboard.service
Wants=chessboard.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/Desktop/chess-board
Environment=PATH=/home/pi/Desktop/chess-board/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/pi/Desktop/chess-board/.venv/bin/python3 /home/pi/Desktop/chess-board/dpad_keyboard.py --mode http --url http://127.0.0.1:8000/api/input
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
